# docker/hbase/hbase.conf
[supervisord]
nodaemon=true

[program:hbase]
command=sh /opt/run.sh
stopwaitsecs=60
