# docker/hbase/hbase.conf
[supervisord]
nodaemon=true

[program:hbase]
command=/opt/run.sh
stopwaitsecs=60
