# device-proxy-service/supervisor-proxyservice.conf
[supervisord]
nodaemon=true

[program:proxyservice]
command=/root/semiot-platform/device-proxy-service/run.sh
stopwaitsecs=20